# deploy/supervisor/caddy.conf
[program:caddy]
command=/usr/bin/caddy run --config /opt/orchestra/caddy/Caddyfile
directory=/opt/orchestra/caddy
autostart=true
autorestart=true
stdout_logfile=/var/log/caddy/stdout.log
stderr_logfile=/var/log/caddy/stderr.log
