# deploy/supervisor/mcp-server.conf
[program:orchestra-mcp]
command=/opt/orchestra/mcp-server/orchestra-mcp
directory=/opt/orchestra/mcp-server
environment=PORT="3001"
autostart=true
autorestart=true
stdout_logfile=/var/log/orchestra/mcp.log
stderr_logfile=/var/log/orchestra/mcp-error.log
